# embulk/crontab
0 * * * * /scripts/run.sh >> /var/log/cron.log 2>&1
